# wg-ondemand.service
[Unit]
Description=WireGuard On-Demand VPN Daemon
Documentation=https://github.com/vly/wg-ondemand
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
# Setup TC qdisc before starting daemon
ExecStartPre=/usr/local/bin/wg-ondemand-setup-tc /etc/wg-ondemand/config.toml
ExecStart=/usr/local/bin/wg-ondemand -c /etc/wg-ondemand/config.toml
Restart=on-failure
RestartSec=5s

# User and group (run as root for BPF and network management capabilities)
User=root
Group=root

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false

# Required capabilities for eBPF and network management
AmbientCapabilities=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
TasksMax=256

[Install]
WantedBy=multi-user.target
